# pwsmqttdispatcher/rootfs/etc/services.d/pwsmqttdispatcher/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the pwsmqttdispatcher service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare mqtt_topic

export MQTT_HOST=$(bashio::services mqtt "host")
export MQTT_USER=$(bashio::services mqtt "username")
export MQTT_PASSWORD=$(bashio::services mqtt "password")
export MQTT_TOPIC="$(bashio::config 'mqtt_topic')"
export MQTT_CLIENT_ID="$(bashio::config 'mqtt_client_id')"
export PWS_IP="$(bashio::config 'pws_ip')"
export FETCH_INTERVAL="$(bashio::config 'fetch_interval')"
export DEBUG_ENABLED="$(bashio::config 'debug_enabled')"

bashio::log.info "Starting pwsmqttdispatcher"

exec /usr/bin/pwsmqttdispatcher
